# deploy/systemd/epoxiron-web.service
[Unit]
Description=Epoxiron Web
After=network.target docker.service epoxiron-api.service
Requires=docker.service epoxiron-api.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/epoxiron
ExecStart=/usr/bin/docker compose -f deploy/docker-compose.vps.yml up -d web
ExecStop=/usr/bin/docker compose -f deploy/docker-compose.vps.yml stop web
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target

